# netextender/supervisord.conf
[supervisord]
nodaemon=true

[program:netextender]
command=/usr/bin/netExtenderGui
autorestart=true
startretries=0
