# remmina/supervisord.conf
[supervisord]
nodaemon=true

[program:netextender]
command=/usr/bin/netExtenderGui
autorestart=true
startretries=0

[program:remmina]
command=/usr/bin/remmina
autorestart=false
startretries=0
